# pppp.conf
[program:pppp]
command=python3 /home/curator2/pppp/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/pppp/test.err.log
stdout_logfile=/home/curator2/pppp/test.out.log